# profile.d/20oos_oraclexe.sh
source ORACLE_HOME/bin/oracle_env.sh
alias sqlplus='rlwrap sqlplus'
alias rman='rlwrap rman'
