# profile.d/30oos_nodejs.sh
export NODE_PATH=/usr/lib/node_modules
export LD_LIBRARY_PATH=${ORACLE_HOME}/lib:${LD_LIBRARY_PATH}
