# init.d/node4ords.conf
# Upstart config file for node4ords

description "Node4Ords - proxy to APEX using ORDS on Tomcat"

start on (local-filesystems and net-device-up IFACE!=lo)
exec /usr/local/bin/node4ords
